# bin/crontab
# crontab:
# Timed tasks for PledgeBank.
#
# Copyright (c) 2005 UK Citizens Online Democracy. All rights reserved.
# Email: francis@mysociety.org. WWW: http://www.mysociety.org/
#
# $Id: crontab,v 1.4 2005-06-08 10:13:34 francis Exp $

0-59/2 * * * * /data/vhost/www.pledgebank.com/mysociety/pb/bin/frequentupdate
33 * * * * /data/vhost/www.pledgebank.com/mysociety/pb/bin/find-pledge-connections

